# systemd/borgwrapper-verify@.timer
[Unit]
Description=monthly borgwrapper backup verification timer

[Timer]
OnCalendar=monthly
RandomizedDelaySec=14d
Persistent=true

[Install]
WantedBy=multi-user.target
